# tips/hdmi_display.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/etienne/work/setup/display/hdmi_display_handle"
